# fundamentals/file/cron_file/cron/file.cron
# file operation crontab
* * * * * (sleep 1; cd /srv/ && ./main)