# advanced/email/imap/cron/imap.cron
# imap client crontab
* * * * * (sleep 5; ./go/src/imap)