# file/cron_file/cron/file.cron
# file operation crontab
* * * * * (sleep 5; /go/src/cron_file/main)